# deploy/colocacion_diaria.cron
# Cron para Linux

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

58 4 * * * cd /home/sa/colocacion_diaria && PYTHON_BIN=/home/sa/colocacion_diaria/.venv/bin/python bash ./scripts/run_daily_reports.sh >> /home/sa/colocacion_diaria/logs/cron.log 2>&1

